# dialectic/deploy/dialectic.service
[Unit]
Description=Dialectic backend (FastAPI + dialogue engine)
Documentation=file:///root/DwoodAmo/dialectic/CLAUDE.md
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DwoodAmo/dialectic
EnvironmentFile=/root/DwoodAmo/dialectic/.env
ExecStart=/usr/bin/python3 run.py
# Restart on crash AND on kill; 5s backoff to avoid tight loops on broken configs
Restart=always
RestartSec=5
# Graceful shutdown window — lets in-flight WebSocket sessions drain before SIGKILL
TimeoutStopSec=15
# Structured logging to journalctl -u dialectic
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dialectic
# Leave a crash breadcrumb if the loop flaps
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
